# oe-core/99-serial-logger.rules
SUBSYSTEM=="tty", ACTION=="add", KERNEL=="ttyACM0", RUN+="/usr/bin/sudo /usr/bin/restart_gpsd.sh"
